# unify.service
[Unit]
Description=Unify Next.js App
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/unify
Environment=NODE_ENV=production
EnvironmentFile=/srv/unify/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
